# supervisor/openrvdas_logger_manager_sqlite.ini
[program:logger_manager_sqlite]
command=OPENRVDAS_ROOT/venv/bin/python server/logger_manager.py --database sqlite --data_server_websocket :8766 -v -V --no-console
environment=PATH="OPENRVDAS_ROOT/venv/bin:/usr/bin:/usr/local/bin"
directory=OPENRVDAS_ROOT
autostart=false
autorestart=false
startretries=3
killasgroup=true
stderr_logfile=/var/log/openrvdas/logger_manager.stderr
stderr_logfile_maxbytes=10000000 ; 10M
stderr_logfile_maxbackups=100
user=RVDAS_USER
